# deploy/dgcomp.timer
[Unit]
Description=Run DGCOMP-says every 2 hours
Requires=dgcomp.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=2h
# Catch up if the box was off when a tick was due.
Persistent=true
# Add up to 5 min of jitter so we don't hammer the EC API on the hour.
RandomizedDelaySec=5min
Unit=dgcomp.service

[Install]
WantedBy=timers.target
